# initsys/bugyou-cntrl.service
[Unit]
Description=Generic bugyou controller
After=network.target redis.service fedmsg-hub.service
Requires=network.target redis.service fedmsg-hub.service
Documentation=http://bugyou.readthedocs.org/

[Service]
ExecStart=/usr/bin/bugyou-cntrl
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target
